# units/kube-flanneld.service
[Unit]
Description=Flanneld loader
Requires=setup-network-environment.service etcd2.service

[X-Fleet]
Global=true